# server/mcp-server.service
[Unit]
Description=滴个龟龟 MCP Server (AI Agent API)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/digeguigui
ExecStart=/usr/bin/python3 /home/ubuntu/digeguigui/scripts/mcp_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
